# systemd/chamado-health.timer
[Unit]
Description=Run Chamado Health check every minute

[Timer]
OnCalendar=minutely
Persistent=true

[Install]
WantedBy=timers.target
