# systemd/chamado-health.timer
[Unit]
Description=Run Chamado Health check every 2 minutes

[Timer]
OnCalendar=*:0/2
Persistent=true

[Install]
WantedBy=timers.target
